# deploy/vps/burgerbus.service
[Unit]
Description=BurgerBUS Online Order App
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/var/www/BurgerBUS-ONLINE-ORDER
EnvironmentFile=/etc/burgerbus/burgerbus.env
ExecStart=/usr/bin/node src/server.js
Restart=always
RestartSec=5
LogsDirectory=burgerbus
StandardOutput=append:/var/log/burgerbus/app.log
StandardError=append:/var/log/burgerbus/error.log

[Install]
WantedBy=multi-user.target
